# deploy/find-hub-tracker.service
[Unit]
Description=Find Hub Device Tracker
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=tracker
WorkingDirectory=/opt/find-hub-tracker
EnvironmentFile=/opt/find-hub-tracker/.env
ExecStart=/opt/find-hub-tracker/.venv/bin/python -m find_hub_tracker start
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
