# deploy/find-hub-tracker.service
[Unit]
Description=Find Hub Device Tracker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/find-hub-tracker
ExecStart=%h/.local/bin/uv run find-hub-tracker start
Restart=always
RestartSec=30
Environment=PATH=%h/.local/bin:%h/.cargo/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
